# config/supervisor/01_camera.conf
[program:camera]
command=robocam_camera
autostart=true
autorestart=true
user=root
