# config/supervisor/02_arduino.conf
[program:arduino]
command=robocam_arduino
autostart=true
autorestart=true
user=root
