# config/supervisor/03_arm.conf
[program:arm]
command=robocam_arm
autostart=true
autorestart=true
user=root
