# config/supervisor/04_distance.conf
[program:distance]
command=robocam_distance
autostart=false
autorestart=false
user=root
