# config/supervisor/04_distance.conf
[program:distance]
command=robocam_distance
autostart=true
autorestart=true
user=root
