# config/supervisor/05_web.conf
[program:web]
command=robocam_web
autostart=true
autorestart=true
user=root
