# config/supervisor/06_fifo.conf
[program:fifo]
command=robocam_fifo
autostart=true
autorestart=true
user=root
